# dotfiles/config/autostart/cerebro.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Cerebro
Exec=cerebro
StartupNotify=false
Terminal=false
